# scripts/launchd/com.darttext.usfinancials.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.darttext.usfinancials</string>

    <!-- US 재무제표(yfinance) 크롤러. 매주 토요일 04:00 —
         .omc/specs/brainstorming-us-market-data-plane.md Round 7 결정.
         FnGuide(03:00)와 겹치지 않게 1시간 뒤로 배치.
         주의: us_company가 아직 비어있음(investing.com 페이지네이션 미해결로
         유니버스 크롤러 미완성) — 채워지기 전까지 이 작업은 매번 0건 처리로
         조용히 끝난다(company 목록이 비었을 뿐 오류는 아님). 유니버스 완성 후
         자동으로 의미 있는 데이터가 쌓인다. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-s</string>
        <string>/Users/gyuyeong/projects/.venv/bin/python3</string>
        <string>/Users/gyuyeong/projects/dart-text2sql-wiki/scripts/run_us_financials.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/gyuyeong/projects/dart-text2sql-wiki</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>6</integer>
        <key>Hour</key><integer>4</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/gyuyeong/projects/dart-text2sql-wiki/data/logs/us_financials.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gyuyeong/projects/dart-text2sql-wiki/data/logs/us_financials.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/gyuyeong/projects/.venv/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/gyuyeong</string>
    </dict>
</dict>
</plist>
