# scripts/launchd/com.darttext.ususerve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.darttext.ususerve</string>

    <!-- US 티커/섹터 유니버스(investing.com Playwright) 크롤러. 매주 토요일 03:30 —
         us_financials(04:00)보다 먼저 돌려 us_company를 채운다.
         주의(2026-07-12 실측): investing.com 2페이지 이상은 InvestingPro 유료
         게이트라 게스트는 시총 상위 30종목(1페이지)만 수집된다. Pro 접근이 생기면
         us_universe._fetch_screener_page 만 교체하면 337페이지 전체 순회가 된다.
         사전 조건: /Users/gyuyeong/projects/.venv 에 playwright 설치 +
         `playwright install chromium`(자체 크로미움 다운로드) 완료. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-s</string>
        <string>/Users/gyuyeong/projects/.venv/bin/python3</string>
        <string>/Users/gyuyeong/projects/dart-text2sql-wiki/scripts/run_us_universe.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/gyuyeong/projects/dart-text2sql-wiki</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>6</integer>
        <key>Hour</key><integer>3</integer>
        <key>Minute</key><integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/gyuyeong/projects/dart-text2sql-wiki/data/logs/us_universe.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gyuyeong/projects/dart-text2sql-wiki/data/logs/us_universe.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/gyuyeong/projects/.venv/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/gyuyeong</string>
    </dict>
</dict>
</plist>
