# scripts/launchd/com.darttext.usfinancialssec.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.darttext.usfinancialssec</string>

    <!-- SEC companyfacts(무료 공식 XBRL) 주간 갱신. 매주 토요일 05:00 —
         .omc/specs/brainstorming-sec-edgar-us-financials-backfill.md AC5.
         기존 배치와 겹치지 않게 배치: FnGuide(03:00) / us_financials yfinance(04:00) 뒤 05:00.
         종목별 SEC API(초당 10건 이하, User-Agent 헤더)로 최신 분기만 upsert 한다.
         초기 대량 백필(companyfacts.zip)은 이 스케줄이 아니라 사용자 감독 하에
         scripts/backfill_us_financials_sec.py 로 1회 별도 실행한다. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-s</string>
        <string>/Users/gyuyeong/projects/.venv/bin/python3</string>
        <string>/Users/gyuyeong/projects/dart-text2sql-wiki/scripts/run_us_financials_sec.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/gyuyeong/projects/dart-text2sql-wiki</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>6</integer>
        <key>Hour</key><integer>5</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/gyuyeong/projects/dart-text2sql-wiki/data/logs/us_financials_sec.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gyuyeong/projects/dart-text2sql-wiki/data/logs/us_financials_sec.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/gyuyeong/projects/.venv/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/gyuyeong</string>
    </dict>
</dict>
</plist>
