# scripts/launchd/com.darttext.usprices.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.darttext.usprices</string>

    <!-- US 주가(yfinance) 크롤러. 매일 08:00 — 미국 정규장 마감(동절기 05:00,
         하절기 06:00 KST) 이후 여유를 두고 실행. usfinancials(토 04:00)와
         달리 매일 실행하는 이유는 주가는 일 단위 갱신이 필요하기 때문
         (증분 수집 기본 구간 5일이라 하루 이틀 밀려도 다음 실행에서 자동 보정). -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-s</string>
        <string>/Users/gyuyeong/projects/.venv/bin/python3</string>
        <string>/Users/gyuyeong/projects/dart-text2sql-wiki/scripts/run_us_prices.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/gyuyeong/projects/dart-text2sql-wiki</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>8</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/gyuyeong/projects/dart-text2sql-wiki/data/logs/us_prices.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gyuyeong/projects/dart-text2sql-wiki/data/logs/us_prices.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/gyuyeong/projects/.venv/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/gyuyeong</string>
    </dict>
</dict>
</plist>
